# negotiator/negotiator-agent/rpm/lib/systemd/system/negotiator-agent.service
[Unit]
Description=negotiator-agent

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/usr/sbin/negotiator-agent
Restart=always
