# negotiator/negotiator-agent/rpm/postinst
#!/bin/bash -e

systemctl daemon-reload
systemctl enable negotiator-agent
systemctl start negotiator-agent
