# negotiator/negotiator-agent/rpm/prerm
#!/bin/bash -e

systemctl stop negotiator-agent
systemctl disable negotiator-agent
