# guest_additions/rpm/lib/systemd/system/testo-guest-additions.service
[Unit]
Description=testo-guest-additions

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/usr/sbin/testo-guest-additions
Restart=always
