# negotiator/negotiator-agent/deb/postinst
#!/bin/bash -e

update-rc.d negotiator-agent defaults
service negotiator-agent start
