# negotiator/negotiator-agent/deb/prerm
#!/bin/bash -e

service negotiator-agent stop
update-rc.d negotiator-agent remove
