# guest_additions/pkg/postinst
#!/bin/bash -e

systemctl daemon-reload
systemctl enable testo-guest-additions
systemctl start testo-guest-additions
