# guest_additions/pkg/prerm
#!/bin/bash -e

systemctl stop testo-guest-additions
systemctl disable testo-guest-additions
