# guest_additions/deb/prerm
#!/bin/bash -e

service testo-guest-additions stop
update-rc.d testo-guest-additions remove
